# profile.d/vbox.sh
#!/bin/sh
export VBOX_PROGRAM_PATH=/usr/lib64/virtualbox/
export VBOX_SDK_PATH=/usr/share/virtualbox/sdk/
