# profile.d/vbox.csh
#!/bin/csh
setenv VBOX_PROGRAM_PATH /usr/lib64/virtualbox/
setenv VBOX_SDK_PATH /usr/share/virtualbox/sdk/
